# scripts/systemd/ns-bootstrap-update-daily.service
[Unit]
Description=ns-bootstrap daily background update (apt)

[Service]
Type=oneshot
ExecStart=%h/ns-bootstrap/scripts/scheduled-update-daily.sh
