# scripts/systemd/ns-bootstrap-update-interactive.timer
[Unit]
Description=Run ns-bootstrap daily interactive update

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
