# scripts/systemd/ns-bootstrap-update-daily.timer
[Unit]
Description=Run ns-bootstrap daily background update

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
